# systemd/hermes-gateway.path
[Unit]
Description=Watch for Hermes config.yaml and start the gateway
Documentation=https://github.com/NousResearch/hermes-agent

[Path]
# Start hermes-gateway.service the moment config.yaml is created.
# After first boot, running `hermes setup` creates config.yaml and the gateway comes up immediately.
# A .path unit automatically triggers the .service with the same name (hermes-gateway).
PathExists=/var/lib/hermes/config.yaml

[Install]
WantedBy=multi-user.target
