# systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway
Documentation=https://github.com/NousResearch/hermes-agent
After=network-online.target
Wants=network-online.target
# Skip quietly when config.yaml is absent (avoids a crash loop).
# Right after first boot there is no config yet, since the user has not run `hermes setup`.
ConditionPathExists=/var/lib/hermes/config.yaml

[Service]
Type=simple
# Data (config.yaml, .env, sessions, logs, ...) lives under /var/lib/hermes.
Environment=HERMES_HOME=/var/lib/hermes
# Gateway start command.
# Reference: install.sh maybe_start_gateway() runs `hermes gateway install` then
#   `hermes gateway start` under systemd, and `nohup hermes gateway` otherwise.
#   So the direct foreground command is `hermes gateway`.
#   - https://raw.githubusercontent.com/NousResearch/hermes-agent/main/scripts/install.sh
# With systemd Type=simple, ExecStart runs `hermes gateway` directly in the foreground.
ExecStart=/usr/local/bin/hermes gateway
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
